# rootfs/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::sysinit:/sbin/mdev -s
::respawn:-/bin/sh
::restart:/sbin/init
::ctrlaltdel:/bin/umount -a -r
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff –a


